# deploy/tentacle-backend.service
[Unit]
Description=Tentacle Backend API
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/tentacle-backend
EnvironmentFile=/opt/tentacle-backend/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tentacle-backend

# Securite
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/tentacle-backend

[Install]
WantedBy=multi-user.target
