# deploy/tentacle-backend.service
[Unit]
Description=Tentacle Backend API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/repo/tentacle-source/apps/backend
EnvironmentFile=/var/repo/tentacle-source/apps/backend/.env
ExecStart=/usr/bin/node dist/index.js
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tentacle-backend
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/repo/tentacle-source

[Install]
WantedBy=multi-user.target
